# hosts
#!/usr/bin/perl

$| = 1;

my $lsb = q{### BEGIN INIT INFO
# Provides:       hosts
# Required-Start: $network $remote_fs $syslog
# Required-Stop:  $network $remote_fs $syslog
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Description:    Host Availability Script
### END INIT INFO
};

use strict;
use warnings;
use FindBin;
use lib "$FindBin::Bin/.";

use Term::ANSIColor qw(:constants :pushpop);
use Term::ANSIColor;

sub init {
	    mkdir "/etc/hosts.d";
}

sub rebuild {
	my $hosts = "";

	if (-f "/etc/hosts.local") {
		open my $f, "<", "/etc/hosts.local";
		my $line;
		while (defined($line = <$f>)) {
			$hosts .= $line;
		}
		close $f;
	}

	$hosts =~ s/^\s+//igm;
	$hosts =~ s/\s+$//igm;

	open my $f, ">", "/etc/hosts" or return; # exit if not root!
	print $f $hosts."\n\n" if $hosts;
	
	my $hostname = hostname();
	
	print $f qq{# static hosts IPv4
127.0.0.1 localhost localhost.localdomain
127.0.1.1 $hostname

# static hosts IPv6
::1 ip6-localhost ip6-loopback
fe00::0 ip6-localnet
ff00::0 ip6-mcastprefix
ff02::1 ip6-allnodes
ff02::2 ip6-allrouters
ff02::3 ip6-allhosts};
	
	print $f "\n\n# hosts.d\n";

    my %hosts_lan = ip_addresses("lan");
	for my $host (keys %hosts_lan) {
		my $ip = $hosts_lan{$host};

		print $f "$ip $host $host.lan\n";
	}
    my %hosts_wan = ip_addresses("wan");
	for my $host (keys %hosts_wan) {
		my $ip = $hosts_wan{$host};

		if ($hosts_lan{$host}) {
			print $f "$ip $host.wan\n";
		}
		else {
			print $f "$ip $host $host.wan\n";
		}
	}

	close $f;
}

sub hostname {
	open my $f, "<", "/etc/hostname";
	my $hostname = <$f>;
	chomp $hostname;
	close $f;
	return $hostname;
}

sub fqdn {
	open my $f, "<", "/etc/mailname";
	my $hostname = <$f>;
	chomp $hostname;
	close $f;
	return $hostname;
}

sub list {
    my %hosts = ip_addresses("lan");
    for my $host (sort keys %hosts) {
        my $ip = $hosts{$host};
        print "$ip $host\n";
    }
}

sub padding {
	my ($str, $len1) = @_;
	my $len2 = length($str);
	if ($len2 <= $len1) {
		return " " x ($len1-$len2);
	}
	else {
		return "";
	}
}

sub status {
    my %hosts_lan = ip_addresses("lan");
    my %hosts_wan = ip_addresses("wan");
    my @hostnames = sort { ([split(/\./, $hosts_lan{$b}||"")]->[0]||0) <=> ([split(/\./, $hosts_lan{$a}||"")]->[0]||0) }
                    sort { ([split(/\./, $hosts_lan{$a}||"")]->[3]||0) <=> ([split(/\./, $hosts_lan{$b}||"")]->[3]||0) }
                    keys %{{%hosts_lan, %hosts_wan}};

    my %ping_ms = ping_ms(map {grep {$_} ($hosts_lan{$_}, $hosts_wan{$_})} @hostnames);
	my %services_hash_lan = services("lan");
	my %services_hash_wan = services("wan");

    for my $host (@hostnames) {
        my $ip_lan = $hosts_lan{$host};
        my $ip_wan = $hosts_wan{$host};

		my $show_host = 1;
	    if ($ip_lan) {
	        print_status_line($ip_lan, $host, "lan", 3, ".lan", \%ping_ms, \%services_hash_lan, $show_host);
	        $show_host = 0;
	    }
	    if ($ip_wan) {
	        print_status_line($ip_wan, $host, "wan", 3, ".wan", \%ping_ms, \%services_hash_wan, $show_host);
	        $show_host = 0;
	    }
    }
}

sub colored_ip {
	my ($ip) = @_;
	
	if ($ip =~ /192.168./) {
		if ($ip =~ /192.168.100./) {
			return "rgb410";
		}
		else {
			return "rgb040";
		}
	}
	else {
		return "cyan";
	}
}

sub print_status_line {
	my ($ip, $host, $network, $network_width, $suffix, $ping_ms_ref, $services_hash_ref, $show_host) = @_;
	
    local $Term::ANSIColor::AUTORESET = 1;
    if ($show_host) {
	    print color("rgb015").$host.RESET;
	    print padding($host, 11);
	}
	else {
		print padding("", 11);
	}
    print " ";
    my $color = colored_ip($ip);
#    print color($color).$ip.RESET;
    print $ip;
    print padding($ip, 15);
    print " (";
    print color($color).$network.RESET;
    print "): ";
    print padding($network, $network_width);

    my $ping = $ping_ms_ref->{$ip};
    
    my @services = $> ? (keys %{$services_hash_ref->{$ip}}) : check_ports($ip);

    if ($ping || @services) {
     	# print BOLD GREEN "up";
    	if ($ping) {
	     	print GREEN "up";
	       	print MAGENTA " ".$ping." ms";
	        print padding($ping, 5);
	    }
	    else {
	     	print RED "unreachable";
	    }

#        my @services = $> ? (keys %{$services_hash_ref->{$ip}}) : check_ports($ip);

        if (@services) {
        	my @sorted_services = sort_services_ports(sort_services_ssl(@services));
           	
            print " ";
            print join(", ", sort @sorted_services);
            set_services($suffix, $host, @services);
        }
    }
    else {
      	print RED "down";
    }
    print "\n";
}

sub sort_services_ports {
	my @services = @_;
	my %ports = ();
	for my $s (@services) {
		my $name = "";
		my $port = 0;
		if ($s =~ /^(.*)[:]([0-9]+)$/) {
			$name = $1;
			$port = $2;
		}
		else {
			$name = $s;
			$port = 1;
		}
		$ports{$name} = [] if !$ports{$name};
		push @{$ports{$name}}, $port;
	}
	@services = ();
	for my $name (keys %ports) {
		my @ps = @{$ports{$name}||[]};
		if (grep { $_ == 1 } @ps) {
			@ps = grep { $_ != 1 } @ps;
			push @services, $name;
		}
		if (@ps > 0) {
			if (@ps > 1) {
				push @services, $name.":[".join(",", @ps)."]";
			}
			else {
				push @services, $name.":".join(",", @ps);
			}
		}
	}
	return @services;
}

sub sort_services_ssl {
	my @services = @_;
    my %services_tmp = map {$_ => 1} @services;
    for my $s (keys %services_tmp) {
    	if ($services_tmp{$s."s"}) {
    		delete $services_tmp{$s."s"};
    		delete $services_tmp{$s};
    		$services_tmp{$s."[s]"} = 1;
        }
    }
    return keys %services_tmp;
}

sub ping_ms {
	my (@ips) = @_;
	my $ips_arg = join(" ", @ips);
	my $output = $ips_arg ? `fping -c1 -t1200 $ips_arg 2>/dev/null` : "";
	my %ping_ms = ();
	for my $line (split(/[\r\n]+/, $output)) {
		if ($line =~ /^([^\s]+?)\s.*?\s([^\s]+?)\sms\s/) {
			$ping_ms{$1} = $2;
		}
	}
	return %ping_ms;
}

sub check_ports {
	my ($ip) = @_;
	my %ports = (
			53 => "dns", 21 => "ftp", 445 => "samba", 515 => "lpr:515", 9100 => "jetdirect",
			25 => "smtp", 465 => "smtps", 143 => "imap",993 => "imaps", 
			22 => "ssh:22", 2200 => "ssh:2200", 2201 => "ssh:2201", 2202 => "ssh:2202", 2203 => "ssh:2203",
			80 => "http", 443 => "https", 8080 => "http:8080", 38080 => "http:38080", 8001 => "lighttpd",
#			8081 => "tail", 38081 => "tail:38081",
			25565 => "minecraft", 10011 => "teamspeak", 5060 => "sip:5060", 5061 => "sip:5061", 64738 => "mumble", 3306 => "mysql",
			17500 => "dropbox", 48888 => "torrentsync:48888", 48887 => "torrentsync:48887", 8888 => "torrentsync-ui",
			5900 => "vnc:0", 5901 => "vnc:1", 5902 => "vnc:2", 5903 => "vnc:3", 5904 => "vnc:4", 5905 => "vnc:5",
	);
	
	my $portlist = join(",", keys %ports);
	my $output = `nmap --host-timeout=500 -sS -oG - -p T:$portlist $ip 2>&1 | egrep '(Ports:)|(root privileges)'`;
	my @running = ();

	if ($output =~ /root privileges/) {
		#for my $port (keys %ports) {
		#	my $service = $ports{$port};
		#	if (`nc -z $ip $port && echo true` =~ /true/) {
		#		push @running, $service;
		#	}
		#}
	}
	else {
		$output =~ s/^(.*?)Ports..//igm;
		my @scanresult = map { s/\s+//igm; $_ } split(/,/, $output);
		@running = map { $ports{$_} } map { [split(/\//, $_)]->[0] } grep {/open/} @scanresult;
	}
	
	return @running;
}

sub get_suffix {
	my ($suffix) = shift @_;
	$suffix = "" if !$suffix;
	$suffix = ".".$suffix if $suffix && $suffix !~ /\./;
	return ($suffix, @_);
}

sub ip_addresses {
	(my $suffix, @_) = get_suffix(@_);

    system("bash -c 'rm -f /etc/hosts.d/*/.??* 2>/dev/null'");
	my %hosts = ();
    for my $host (split /\s+/, `cd /etc/hosts.d/; ls`) {
        chomp $host;

        # skip if host has been deleted
        next if -f "/etc/hosts.d/".$host."/remove";

		# hack! for migration to new filenames...
		unlink "/etc/hosts.d/".$host."/ip";
		unlink "/etc/hosts.d/".$host."/ip.real";
		unlink "/etc/hosts.d/".$host."/services";
		
        open my $h, "<", "/etc/hosts.d/".$host."/ip".$suffix or next;
        my $ip = <$h>;
        chomp $ip;
        close $h;
        
        if ($ip =~ /[a-zA-Z]/) {
        	$ip = `gethostip -d "$ip"`;
        }
       	$ip =~ s/[^0-9.]//igm;

        $hosts{$host} = $ip;
    }
    return %hosts;
}

sub services {
	(my $suffix, @_) = get_suffix(@_);

	my %services = ();
	my %ips = ip_addresses($suffix);
    for my $host (keys %ips) {
        chomp $host;
        open my $h, "<", "/etc/hosts.d/".$host."/services".$suffix or next;
        my @s = <$h>;
        close $h;
        @s = split(/[^a-z0-9:]/, join(" ", @s));
        $services{$ips{$host}} = $services{$host} = {map {$_ => 1} @s};
    }
    return %services;
}

sub set_services {
	(my $suffix, @_) = get_suffix(@_);
    my ($host, @services) = @_;

    open my $h, ">", "/etc/hosts.d/".$host."/services".$suffix or return;
    print $h join(" ", @services)."\n";
    close $h;
}

sub add {
	my @arg = @_;
	my $suffix = "";

	my %is_possible_network = (lan => 1, wan => 1);
	if (@arg > 0 && $is_possible_network{$arg[0]}) {
		$suffix = "." . $arg[0];
		shift @arg;
	}

	if (scalar @arg == 2) {
	    my $ip = $arg[0];
	    my $host = $arg[1];
	    $ip =~ s/[^a-zA-Z0-9_.-]//igm;
	    $host =~ s/[^a-zA-Z0-9_.-]//igm;
	    if ($ip =~ /[a-zA-Z]/i) {
	    	($host, $ip) = ($ip, $host);
	    }
	    
	    print "IP: " . $ip . "\n";
	    print "Host: " . $host . "\n";
	    print "Network: " . substr($suffix, 1) . "\n";
	    
	    mkdir "/etc/hosts.d/".$host;
	    open my $f, ">", "/etc/hosts.d/".$host."/ip".$suffix;
	    print $f $ip."\n";
	    close $f;
	}
	else {
		print STDERR "Syntax: $0 add [lan|wan] IP HOST\n";
	}
}

sub sync_dirs {
	my %syncdirs = ();
	my %ips = ip_addresses("lan");
    for my $host (keys %ips) {
        chomp $host;
        open my $h, "<", "/etc/hosts.d/".$host."/sync" or next;
        my %dirs = ();
        for my $line (<$h>) {
        	if ($line =~ /^([^\s]+?)\s*?[=]\s*(.*)\s*/) {
        		$dirs{$1} = $2;
        	}
        }
        close $h;
        $syncdirs{$ips{$host}} = $syncdirs{$host} = {%dirs};
    }
    return %syncdirs;
}

sub install {
	# do not install anything when run from /etc/init.d !!
	if ($0 =~ /init.d/) {
		return;
	}
		
	{
		system("mv /etc/crontab /etc/crontab.bak; grep -v $0 /etc/crontab.bak | grep -v HOSTS > /etc/crontab");
	}
	{
		system("rm /etc/resolv.conf");
		open my $f, ">", "/etc/resolv.conf";
		print $f qq{nameserver 8.8.8.8\n};
		close $f;
	}
}

init;
if (@ARGV && $ARGV[0] eq "add") {
	shift @ARGV;
	add @ARGV;
}
elsif (@ARGV && $ARGV[0] eq "sync") {
}
elsif (@ARGV && $ARGV[0] eq "list") {
	list;
}
elsif (@ARGV && $ARGV[0] eq "status") {
	status;
}
elsif (@ARGV && ($ARGV[0] eq "start" || $ARGV[0] eq "restart" || $ARGV[0] eq "install")) {
    print "Starting hosts script...\n" if $ARGV[0] eq "start";
    print "Restarting hosts script...\n" if $ARGV[0] eq "restart";
	install;
	rebuild;
}
elsif (@ARGV && ($ARGV[0] eq "stop")) {
    print "Stopping hosts script...\n";
    rebuild;
}
else {
	rebuild;
}

1;
