# tinc-extras/etc/network/if-up.d/tinc
#!/bin/sh

set -e

[ "$METHOD" = loopback ] && exit 0
[ -n "$IF_TINC_NET" ] && exit 0

invoke-rc.d tinc alarm || exit 0
